# i2p-supervisor.conf
[program:i2p]
command=/usr/bin/i2prouter console
autostart = true
user=retrouser
startsecs=0
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10

